# deploy/rely/tmp/supervisor/192.168.112.33_15678/node-15678.conf
[program:node-15678]
command=/home/juzhen/./platon_test/node-15678/platon --identity platon --datadir /home/juzhen/./platon_test/node-15678/data --port 15678 --syncmode 'full' --debug --verbosity 4 --pprof --pprofaddr 192.168.112.33 --pprofport 6060 --metrics --metrics.influxdb --metrics.influxdb.username admin --metrics.influxdb.password admin --metrics.influxdb.endpoint http://192.168.10.102:8086 --metrics.influxdb.host.tag node-tag --rpc --rpcaddr 0.0.0.0 --rpcport 5678 --rpcapi platon,debug,personal,admin,net,web3 --txpool.nolocals --networkid 100 --gcmode archive --nodekey /home/juzhen/./platon_test/node-15678/data/nodekey --config /home/juzhen/./platon_test/node-15678/config.json
environment=LD_LIBRARY_PATH=/home/juzhen/mpclib
numprocs=1
autostart=false
startsecs=3
startretries=3
autorestart=unexpected
exitcode=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/home/juzhen/./platon_test/node-15678/log/platon.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=10
