# deploy/rely/tmp/supervisor/10.10.8.16_16786/node-16786.conf
[program:node-16786]
command=/home/juzhen/./platon_test/node-16786/platon --identity platon --datadir /home/juzhen/./platon_test/node-16786/data --port 16786 --syncmode 'full' --debug --verbosity 4 --rpc --rpcaddr 0.0.0.0 --rpcport 6786 --rpcapi platon,debug,personal,admin,net,web3 --txpool.nolocals --gcmode archive --nodekey /home/juzhen/./platon_test/node-16786/data/nodekey
environment=LD_LIBRARY_PATH=/home/juzhen/mpclib
numprocs=1
autostart=false
startsecs=3
startretries=3
autorestart=unexpected
exitcode=0
stopsignal=TERM
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/home/juzhen/./platon_test/node-16786/log/platon.log
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=10
